# docs/systemd/rlens.service
[Unit]
Description=rLens Web UI (rlens user-unit example/template)
After=network.target

# Operating note (TASK-SERVICE-001, 2026-06-03):
# Python lädt Module beim Prozessstart. Nach jedem `git pull`/Merge in %h/repos/lenskit
# MUSS dieser User-Service neu gestartet werden, damit die Dumps gegen den frischen
# Code laufen:
#   systemctl --user restart rlens
# Ein laufender Service sieht neue Repo-Dateien nicht automatisch; ohne Restart
# emittiert er Dumps aus einer eingefrorenen Code-Sicht (z. B. fehlende
# generator.runtime, post_emit_health_path, bundle_surface_validation_path,
# bundle_surface_validation_status, claim_evidence_map_json). Beweis + Smoke-Script:
# docs/proofs/rlens-post-merge-restart-smoke-proof.md
# scripts/rlens-post-merge-surface-smoke.sh

[Service]
Type=simple
Environment=RLENS_HUB=%h/repos
Environment=RLENS_MERGES=%h/lenskit-out
Environment=RLENS_HOST=127.0.0.1
Environment=RLENS_PORT=8787
# Required local file; do not commit secrets.
# RLENS_TOKEN is required for API auth.
# RLENS_FS_TOKEN_SECRET is separate and optional for filesystem-token signing.
# Missing file intentionally fails the unit for full-peer deployments.
EnvironmentFile=%h/.config/rlens/env
ExecStartPre=/bin/sh -c 'test -n "$RLENS_TOKEN" || { echo "RLENS_TOKEN is required" >&2; exit 1; }'
WorkingDirectory=%h/repos/lenskit
ExecStart=/usr/bin/python3 -m merger.lenskit.cli.rlens
Restart=on-failure
RestartSec=2

[Install]
WantedBy=default.target
