# docs/systemd/rlens.service
[Unit]
Description=rLens Web UI (rlens user-unit example/template)
After=network.target

[Service]
Type=simple
Environment=RLENS_HUB=%h/repos
Environment=RLENS_MERGES=%h/lenskit-out
Environment=RLENS_HOST=127.0.0.1
Environment=RLENS_PORT=8787
# Required local file; do not commit secrets.
# RLENS_TOKEN is required for API auth.
# RLENS_FS_TOKEN_SECRET is separate and optional for filesystem-token signing.
# Missing file intentionally fails the unit for full-peer deployments.
EnvironmentFile=%h/.config/rlens/env
ExecStartPre=/bin/sh -c 'test -n "$RLENS_TOKEN" || { echo "RLENS_TOKEN is required" >&2; exit 1; }'
WorkingDirectory=%h/repos/lenskit
ExecStart=/usr/bin/python3 -m merger.lenskit.cli.rlens
Restart=on-failure
RestartSec=2

[Install]
WantedBy=default.target
