# ops/systemd/systemkatalog-publish/repobrief-publish-heimgewebe-katalog-main-watch.service
[Unit]
Description=RepoBrief Heimgewebe-Systemkatalog main watcher

[Service]
Type=oneshot
ExecStart=/home/alex/.local/bin/repobrief-publish-heimgewebe-katalog-main-if-changed
Nice=10
IOSchedulingClass=best-effort
IOSchedulingPriority=6
