# ops/systemd/systemkatalog-publish/repobrief-publish-systemkatalog-main.timer
[Unit]
Description=RepoBrief Systemkatalog fallback timer

[Timer]
OnCalendar=daily
RandomizedDelaySec=15min
Persistent=true
Unit=repobrief-publish-systemkatalog-main.service

[Install]
WantedBy=timers.target
