# ops/systemd/systemkatalog-publish/repobrief-publish-heimgewebe-katalog-main.timer
[Unit]
Description=RepoBrief Heimgewebe-Systemkatalog fallback timer

[Timer]
OnCalendar=daily
RandomizedDelaySec=15min
Persistent=true
Unit=repobrief-publish-heimgewebe-katalog-main.service

[Install]
WantedBy=timers.target
