# ops/systemd/repoground-fleet/repoground-publish-fleet-watch.timer
[Unit]
Description=Hourly changed-only RepoGround fleet publication

[Timer]
OnCalendar=hourly
RandomizedDelaySec=10min
Persistent=true
Unit=repoground-publish-fleet-watch.service

[Install]
WantedBy=timers.target
